# jellyfin-media-player/pulseaudio.service.d/30-match-table.conf
[Service]
# Lower the music volume compared to the video volume
ExecStartPost=/usr/bin/pacmd load-module module-match key=media.role table=/etc/pulse/match.table
