# jellyfin-media-player/pulseaudio.service.d/50-dbus.conf
[Service]
# Dbus support for my custom volume notification daemon,
# and custom role corking synchronised with Snapclient
ExecStartPost=/usr/bin/pacmd load-module module-dbus-protocol
