# jellyfin-media-player/pulseaudio.service.d/90-vol-base.conf
[Service]
# Set volume of combined sink to 150%, which we display as 75% to the users
ExecStartPost=/usr/bin/pacmd set-sink-volume combined 65536
