# jellyfin-media-player/systemd-user/audio-muted.target
# This target indicates that the audio volume has been muted
[Unit]
Description=audio muted status
