# jellyfin-media-player/systemd-user/audio-muted_inverter.service
# I want tasmota-audio.service to power off on stop,
# otherwise whenever I shut the system down it will power *on* the speakers.
# Problem is I need audio-muted.target to match media-playback.target in that it will be started when I want the speakers off.
#
# So this unit should, in theory, start tasmota-audio.service on stop, and stop it on start
# FIXME: Wait, what happens with this service on shutdown?
[Unit]
Description=invert audio-mute.target state for audio-output.target
Before=tasmota-audio.service audio-output.target
PartOf=audio-muted.target

[Service]
Type=oneshot
RemainAfterExit=true

ExecStart=systemctl --user stop --no-block audio-output.target
ExecStop=systemctl --user start --no-block audio-output.target

[Install]
WantedBy=audio-muted.target
