# jellyfin-media-player/systemd-user/video-output.target
# This target indicates that the TV should be on
[Unit]
Description=TV power state
